# deploy/tokilake.service
# Suggested file path: /etc/systemd/system/tokilake.service
# Usage:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tokilake
#   sudo systemctl status tokilake
[Unit]
Description=Tokilake Service
After=network.target

[Service]
Type=simple
# Modify these paths/users before use.
User=tokilake
Group=tokilake
WorkingDirectory=/opt/tokilake
ExecStart=/opt/tokilake/dist/one-api --config /opt/tokilake/dist/config.yaml --port 19981 --log-dir /opt/tokilake/dist/logs
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
